# caps.service
[Unit]
Description=Caps Lock off

[Service]
ExecStart=/usr/bin/Caps_Lock.sh

[Install]
WantedBy=multi-user.target
